# polyserver.service
# systemd unit for the poly-server coordinator (DESIGN.md §3, §8).
#
# Install on the droplet:
#   sudo cp /home/kylea/polyselect-distributed/polyserver.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now polyserver        # start now + on boot
#   systemctl status polyserver                   # check it's running
#   journalctl -u polyserver -f                   # tail logs (replaces the screen session)
#   sudo systemctl stop polyserver                # clean shutdown (no more zombies)
#
# Edit the paths/port/jobdir below if they change. Single process by design — do NOT add
# uvicorn --workers (it would break the shared-connection + lock invariant; DESIGN.md §9).
# Binds 0.0.0.0: firewall the port to your worker IPs (auth is just the worker token).

[Unit]
Description=polyselect-distributed coordinator (poly-server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/kylea/polyselect-distributed
# --lease-seconds MUST exceed the worst-case per-coefficient runtime, or the sweep will
# requeue an in-flight coeff and the client's eventual /submit gets a 409 and is dropped.
# For a C209, per-coeff stage 1 is multi-hour: 86400 = 24h. Tune to ~2-3x measured coeff time.
ExecStart=/home/kylea/polyselect-distributed/.venv/bin/python -m polyserver serve \
    --jobdir /home/kylea/polyselect-distributed/AS552 --bind 0.0.0.0 --port 8084 \
    --lease-seconds 86400
# uvicorn handles SIGTERM gracefully, so systemd's default stop is a clean shutdown.
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
